# etc/systemd/system/cam_youtube.service.d/myenv.conf
[Service]
# Get your stream key from your YouTube Studio 'livestreaming' page.
Environment="STREAM_KEY=put-your-steam-key-here"
